# etc/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

00 * * * * root /usr/bin/env bash -c 'freshclam --quiet' 2>&1
10 * * * * root /usr/bin/env bash -c 'clam-test'